# systemd/redshift_adjust/redshift_adjust.timer
[Unit]
Description=redshift_adjust

[Timer]
OnBootSec=1min
OnCalendar=*:0/2
Unit=redshift_adjust.service

[Install]
WantedBy=timers.target
