# systemd/git_autocommit/git-autocommit@.service
[Unit]
Description=Auto-commit for %I

[Service]
Type=oneshot
ExecStart=%h/bin/git_autocommit "%I"
