# systemd/drk_daemon/drk_daemon.target
[Unit]
BindsTo=dev-drk.device
After=dev-drk.device

[Install]
WantedBy=multi-user.target
