# systemd/dalsik_daemon/dalsik_daemon.target
[Unit]
BindsTo=dev-dalsik.device
After=dev-dalsik.device

[Install]
WantedBy=dev-dalsik.device
