# systemd/clean_daily/clean_daily.service
[Unit]
Description=clean daily

[Service]
Type=oneshot
ExecStart=%h/bin/clean_daily.sh

[Install]
WantedBy=multi-user.target
