# udev/99-ryzenadj-hook.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/home/davs/dotfiles/bin/ryzenadj_preset"
